# cookbooks/pgweb/templates/default/pgweb.init.d.erb
#!/sbin/runscript
# $Id$

NAME="pgweb"
PIDFILE="/var/run/$NAME.pid"
USER="deploy"
SU="su $USER -s /bin/bash"
TIMEOUT=5
description="pgweb"
command="/usr/bin/pgweb"
command_args=" --bind=0.0.0.0 --listen=8089 --bookmarks-dir=/home/deploy/.pgweb/bookmarks/ --bookmark=wellthie --readonly -s"

start() {
  if [ -f $PIDFILE ]; then
      echo "Already running... cat $PIDFILE"
      exit 0
  fi

  /sbin/start-stop-daemon --start -b -- ${command} ${command_args} # Note! Logs are lost.
  sleep 1
  PID=`pidof pgweb`
  if [ -z $PID ]; then
      exit 1
  else
      echo $PID > $PIDFILE
  fi
}

stop() {
  if [ -f $PIDFILE ]; then
      start-stop-daemon --stop --pidfile=$PIDFILE
  else
      kill -9 `pidof pgweb`
  fi

}
